# deploy/systemd/360logger-photo.service
[Unit]
Description=360 photo capture session (mouse-triggered) on boot
After=network-online.target 360logger-boot.service
Wants=network-online.target

[Service]
Type=simple
User=__REPO_USER__
Environment=HOME=/home/__REPO_USER__
Environment=PATH=/home/__REPO_USER__/.local/bin:/usr/local/bin:/usr/bin:/bin
WorkingDirectory=__REPO_DIR__
ExecStart=/usr/bin/python3 -m recorder.photo_session
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
